# run/systemd/generator/-.mount
# Automatically generated by systemd-fstab-generator

[Unit]
SourcePath=/etc/fstab
Documentation=man:fstab(5) man:systemd-fstab-generator(8)
Before=local-fs.target

[Mount]
Where=/
What=/dev/disk/by-uuid/949779ce-46aa-434e-8eb0-852514a5d69e
Type=xfs
